# zigstar-cc2538-flasher/rootfs/etc/services.d/zigstar-cc2538-flasher/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# ZigStar CC2538 Flasher service finished
# ==============================================================================

bashio::log.info "ZigStar CC2538 Flasher service finished"
